# systemd/cash-deploy.service
[Unit]
Description=cash-containers GitOps deploy
After=docker.service network.target
Requires=docker.service

[Service]
Type=oneshot
User=mjb
ExecStart=/home/mjb/projects/cash-compose/deploy.sh
StandardOutput=journal
StandardError=journal
